# install/debian/postinst.in
#!/bin/bash -e

. /usr/share/debconf/confmodule

NAME=fah-client
CLIENT_USER=$NAME
CLIENT_GROUP=$NAME
CLIENT_HOME=/var/lib/$NAME
CLIENT_LOGS=/var/log/$NAME
CLIENT_SHARE=/usr/share/$NAME

case "$1" in
  reconfigure|configure)
    # Create group if it does not exist
    if ! getent group | grep -q "^$CLIENT_GROUP:"; then
      echo "Adding group '$CLIENT_GROUP'."
      groupadd --system $CLIENT_GROUP
    fi

    # Create user if it does not exist
    if ! getent passwd | grep -q "^$CLIENT_USER:"; then
      echo "Adding user '$CLIENT_USER'."
      useradd --system -g $CLIENT_GROUP -K UMASK=002 $CLIENT_USER \
              --groups video,render
    fi

    # Create directories
    install -d -m 755 -o $CLIENT_USER -g $CLIENT_GROUP $CLIENT_HOME $CLIENT_LOGS

    # Create a symbolic link to the configuration file
    ln -sf /etc/fahclient/config.xml $CLIENT_HOME/config.xml

    # Install the service
    if [ ! -e /etc/systemd/system/$NAME.service ]; then
      ln -sf $CLIENT_SHARE/$NAME.service /etc/systemd/system/
    fi
    systemctl daemon-reload

    # Start the service
    systemctl start $NAME
    systemctl enable $NAME

    echo
    echo "The Folding@home client is now installed"
    echo
    echo "Access the web interface by going to:"
    echo
    echo "  %(PACKAGE_URL)s"
    echo
    echo "File storage locations:"
    echo "  Logs: /var/log/fah-client"
    echo "  Data: /var/lib/fah-client"
    echo
    echo "Related service commands for Status, Start, Stop, Restart:"
    echo "  systemctl status --no-pager -l fah-client"
    echo "  sudo systemctl start fah-client"
    echo "  sudo systemctl stop fah-client"
    echo "  sudo systemctl restart fah-client"
    ;;
esac
