# install/debian/postinst.in
#!/bin/sh -e

NAME=fah-client
CLIENT_HOME=/var/lib/$NAME
CLIENT_LOGS=/var/log/$NAME

case "$1" in
  reconfigure|configure)
    # Create group if it does not exist
    if ! getent group $NAME >/dev/null; then
      groupadd --force --system $NAME
    fi

    # Create user if it does not exist
    if ! getent passwd $NAME >/dev/null; then
      useradd --system --gid $NAME --shell /usr/sbin/nologin \
        --home-dir $CLIENT_HOME --no-create-home --groups video,render $NAME
    fi

    # Create directories
    install -d -m 755 -o $NAME -g $NAME $CLIENT_HOME $CLIENT_LOGS

    # Create a symbolic link to the configuration file
    ln -sf /etc/fahclient/config.xml $CLIENT_HOME/config.xml

    # Start the service
    systemctl daemon-reload
    systemctl enable $NAME
    systemctl start $NAME

    echo
    echo "The Folding@home client is now installed"
    echo
    echo "File locations:"
    echo
    echo "  Logs: /var/log/$NAME"
    echo "  Data: /var/lib/$NAME"
    echo
    echo "Service commands:"
    echo
    echo "  systemctl status --no-pager -l $NAME"
    echo "  sudo systemctl start $NAME"
    echo "  sudo systemctl stop $NAME"
    echo "  sudo systemctl restart $NAME"
    echo
    echo "Access the web interface by going to:"
    echo
    echo "  %(PACKAGE_URL)s"
    echo
    ;;
esac
